# sandbox/supervisord.conf
; 通信 http 服务配置
[unix_http_server]
file=/tmp/supervisor.sock    ; socket 文件路径
chmod=0770                   ; socket 文件权限，文件所有者和同组用户可读写执行
chown=ubuntu:ubuntu          ; socket 文件所有者和所属组

; supervisord 日志配置
[supervisord]
logfile=/dev/stdout          ; 主日志路径
logfile_maxbytes=0           ; 不限制日志大小
loglevel=info                ; 日志等级
pidfile=/tmp/supervisord.pid ; pid 文件路径
nodaemon=true                ; 前台运行，便于 Docker 捕获日志
minfds=1024                  ; 打开文件描述符的最小值
minprocs=200                 ; 最小进程数
autoshutdown=true            ; 所有服务停止后自动退出

; 启用 supervisor 的远程过程调用接口
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; supervisorctl 连接服务配置
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; 后端服务进程配置
[program:app]
command=bash -lc "PORT=8080 npm run start:local -- %(ENV_APP_ARGS)s"
directory=/sandbox
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME=/home/ubuntu
priority=60

; Chrome 浏览器配置
[program:chrome]
command=chromium \
    --display=:1 \
    --window-size=1280,1080 \
    --start-maximized \
    --no-sandbox \
    --disable-dev-shm-usage \
    --disable-setuid-sandbox \
    --disable-accelerated-2d-canvas \
    --disable-gpu \
    --disable-features=WelcomeExperience,SigninPromo \
    --no-first-run \
    --no-default-browser-check \
    --disable-infobars \
    --test-type \
    --disable-popup-blocking \
    --disable-gpu-sandbox \
    --no-xshm \
    --new-window=false \
    --disable-notifications \
    --disable-extensions \
    --disable-component-extensions-with-background-pages \
    --disable-prompt-on-repost \
    --disable-dialogs \
    --disable-modal-dialogs \
    --disable-web-security \
    --disable-site-isolation-trials \
    --remote-debugging-address=0.0.0.0 \
    --remote-debugging-port=8222 %(ENV_CHROME_ARGS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=20
startretries=3
startsecs=5

; socat 进程配置
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=2

; Xvfb 虚拟显示器配置
[program:xvfb]
command=bash -c "rm -f /tmp/.X1-lock && Xvfb :1 -screen 0 1280x1080x24"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startsecs=2

; x11vnc 服务配置
[program:x11vnc]
command=x11vnc -display :1 -nopw -shared -listen 0.0.0.0 -xkb -forever -rfbport 5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
startsecs=3

; Websockify 配置，将 VNC 转换为 WebSocket
[program:websockify]
command=websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
startsecs=3

; 使用分组统一启动或关闭进程
[group:services]
programs=xvfb,chrome,socat,x11vnc,websockify,app
